# foo.service
[Unit]
Description=Run foo weekly and on boot

[Service]
Type=oneshot
ExecStart=/usr/bin/echo "this is a test"

[Install]
WantedBy=basic.target
